# backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info

[program:nginx]
command=nginx -c /app/config/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
startretries=5
startsecs=2
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:actions]
command=python -m rasa_sdk --actions actions --port 5055
directory=/app/rasa
autostart=true
autorestart=true
startretries=5
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:rasa]
command=/app/scripts/run_rasa.sh
directory=/app/rasa
autostart=true
autorestart=unexpected
startretries=2
startsecs=30
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:proxy]
command=uvicorn hf_proxy:app --host 127.0.0.1 --port 8080 --loop asyncio --http h11
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=25

[program:streamlit]
command=streamlit run /app/ui/streamlit_app.py --server.port 8502 --server.address 127.0.0.1 --server.headless true --server.enableCORS false --server.enableXsrfProtection false --server.baseUrlPath "/ui"
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=25
